# check-ip-cron/cron.d/check_ip
*/5 * * * * /check_ip.sh >> /var/log/cron.log 2>&1